# systemd/voice-trace-sweep.timer
[Unit]
Description=Daily NanoClaw voice-trace retention sweeper
Requires=voice-trace-sweep.service

[Timer]
# Daily at 04:00 local time. Persistent so missed runs (host offline) fire on next boot.
OnCalendar=*-*-* 04:00:00
Persistent=true
RandomizedDelaySec=5m
Unit=voice-trace-sweep.service

[Install]
WantedBy=timers.target
